# collections/ansible_collections/jedimasterjonny/lex/roles/nfs_client/templates/nfs-share.automount.j2
{{ ansible_managed | comment }}
[Unit]
Description=Automount for {{ item.where }}

[Automount]
Where={{ item.where }}
TimeoutIdleSec={{ item.idle_timeout }}

[Install]
WantedBy=multi-user.target
