# roles/podman_backup/templates/podman-backup.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Weekly backup of podman volumes

[Timer]
OnCalendar={{ podman_backup_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
