# roles/homepage/templates/homepage.container.j2
[Unit]
Description=Homepage

[Container]
ContainerName=homepage
Image={{ homepage_image }}
Network=caddy.network
NoNewPrivileges=true
DropCapability=all
{# The root entrypoint chowns the config volume to PUID (CHOWN), then drops to #}
{# it with su-exec (SETUID, SETGID); homepage binds :3000, so no #}
{# NET_BIND_SERVICE. #}
AddCapability=CHOWN SETGID SETUID
Environment=PUID={{ homepage_uid }}
Environment=PGID={{ homepage_uid }}
Environment=TZ={{ homepage_timezone }}
# v1 validates the Host header: the apex vhost, plus the loopback host:port the
# liveness probe sends.
Environment=HOMEPAGE_ALLOWED_HOSTS={{ homepage_domain }},127.0.0.1:3000
# Log to stdout only, so nothing is written under the config mount.
Environment=LOG_TARGETS=stdout
# Operator data store: homepage seeds skeleton services/widgets/bookmarks here on
# first boot and the entrypoint chowns it to PUID. A named volume keeps the
# backing store off a host bind mount — podman labels it container_file_t (no :Z)
# and podman_backup captures it automatically.
Volume=homepage-config:/app/config
# Ansible owns only settings.yaml; mounted read-only as a single file on top of
# the volume so it stays config-as-code without shadowing the seeded data.
Volume=/etc/homepage/settings.yaml:/app/config/settings.yaml:ro,Z
# Liveness against homepage's own health endpoint (busybox wget ships in the
# alpine image). Status only.
HealthCmd=wget -q --tries=1 --spider http://127.0.0.1:3000/api/healthcheck
HealthInterval=1m
HealthTimeout=10s
HealthStartPeriod=30s
HealthRetries=3

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
