# roles/dev/templates/claude-remote-control.service.j2
# {{ ansible_managed }}
[Unit]
Description=Claude Code Remote Control session for {{ dev_user }}
# Server mode polls the Anthropic API and exits on a long network outage; never
# stop retrying so remote access returns on its own once the network does. The
# user manager has no network-online.target to order against, but that retry loop
# also covers a boot start that beats the network up.
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory={{ dev_remote_control_workdir }}
Environment=PATH={{ dev_user_home }}/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
# Remote Control rejects API-key auth, so make sure an inherited key can't break it.
UnsetEnvironment=ANTHROPIC_API_KEY
ExecStart={{ dev_claude_bin }} remote-control
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
