# roles/wordpress/templates/mariadb.container.j2
[Unit]
Description=WordPress database

[Container]
ContainerName=wordpress-db
Image={{ wordpress_db_image }}
Network=caddy.network
NoNewPrivileges=true
DropCapability=all
{# The root entrypoint fixes the datadir's ownership (CHOWN, DAC_OVERRIDE) then #}
{# drops to the mysql user (SETUID, SETGID) before exec'ing mariadbd. #}
AddCapability=CHOWN DAC_OVERRIDE SETGID SETUID
EnvironmentFile=/etc/wordpress/db.env
Volume=wordpress-db:/var/lib/mysql
{# The image's healthcheck.sh: accepts connections and InnoDB is up. #}
HealthCmd=healthcheck.sh --connect --innodb_initialized
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=120s
HealthOnFailure=kill

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
