# roles/wordpress/templates/wordpress.container.j2
[Unit]
Description=WordPress
{# The database unit starts first; wordpress reconnects until it answers. #}
Requires=wordpress-db.service
After=wordpress-db.service
{# The cache is optional — the drop-in degrades gracefully if it is down — so #}
{# only order after it, never require it. #}
Wants=wordpress-redis.service
After=wordpress-redis.service

[Container]
ContainerName=wordpress
Image={{ wordpress_image }}
Network=caddy.network
NoNewPrivileges=true
DropCapability=all
{# The root entrypoint unpacks core into the volume (CHOWN, FOWNER, #}
{# DAC_OVERRIDE), then apache binds :80 (NET_BIND_SERVICE) and drops its #}
{# workers to www-data (SETUID, SETGID). #}
AddCapability=CHOWN DAC_OVERRIDE FOWNER NET_BIND_SERVICE SETGID SETUID
EnvironmentFile=/etc/wordpress/wordpress.env
Volume=wordpress-html:/var/www/html
{# DB-independent liveness — a static asset, so a database blip never #}
{# restart-loops apache; HealthOnFailure=kill + Restart=on-failure self-heals. #}
HealthCmd=curl -fsS http://127.0.0.1/wp-includes/images/blank.gif -o /dev/null
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=120s
HealthOnFailure=kill

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
