# roles/restic_backup/templates/restic-backup.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Scheduled backup of {{ restic_backup_description }}

[Timer]
OnCalendar={{ restic_backup_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
