# roles/podman/files/podman-image-prune.service
[Unit]
Description=Prune dangling podman images

[Service]
Type=oneshot
ExecStart=/usr/local/bin/podman-image-prune
