# roles/podman/files/podman-image-prune.timer
[Unit]
Description=Weekly prune of unused podman images

[Timer]
OnCalendar=Sat *-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
