# roles/podman/files/podman-image-prune.service
[Unit]
Description=Prune unused podman images

[Service]
Type=oneshot
ExecStart=/usr/bin/podman image prune --all --force
