# roles/caddy/templates/caddy.container.j2
[Unit]
Description=Caddy reverse proxy

[Container]
ContainerName=caddy
Image={{ caddy_image }}
Network=caddy.network
NoNewPrivileges=true
DropCapability=all
{# Only to bind the privileged 80/443; caddy needs nothing else. #}
AddCapability=NET_BIND_SERVICE
EnvironmentFile=/etc/caddy/env
PublishPort=80:80
PublishPort=443:443
PublishPort=443:443/udp
Volume=/etc/caddy:/etc/caddy:ro
Volume=caddy-data:/data
Volume=caddy-config:/config
# Liveness probe against caddy's own 204 health endpoint (the localhost site in
# the Caddyfile), so a wedged-but-running proxy is caught. Status only.
HealthCmd=curl -fsS http://localhost/
HealthInterval=1m
HealthTimeout=10s
HealthStartPeriod=30s
HealthRetries=3

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
