# roles/alertmanager/templates/alertmanager.container.j2
[Unit]
Description=Alertmanager

[Container]
ContainerName=alertmanager
Image={{ alertmanager_image }}
# Host network so Prometheus on the NAS reaches the API at :9093 over the LAN,
# like the node and cadvisor exporters. The image's default command already points
# at the mounted config and storage paths. Opening the port is the play's job.
Network=host
NoNewPrivileges=true
DropCapability=all
{# alertmanager runs non-root as nobody (65534), so it needs no capability back. #}
# :Z relabels the mount container_file_t on each start, so it never depends on the
# host's ambient label surviving an SELinux relabel (as /etc/caddy's did not).
Volume=/etc/alertmanager:/etc/alertmanager:ro,Z
# :U hands the named volume to the image's nobody user (65534), which owns the
# silences and notification log under /alertmanager; rootful podman would
# otherwise create it root-owned.
Volume=alertmanager-data:/alertmanager:U
# Restart backstop, not the monitor — see the health convention in CLAUDE.md.
HealthCmd=wget -q -O /dev/null http://localhost:9093/-/healthy
HealthInterval=5m
HealthTimeout=30s
HealthStartPeriod=30s
HealthRetries=3
HealthOnFailure=kill

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
