# roles/gitops_reconcile/templates/gitops-reconcile.service.j2
# {{ ansible_managed }}
[Unit]
Description=Apply origin/{{ gitops_reconcile_branch }} to the fleet
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart={{ gitops_reconcile_script }}
# Runs whether ExecStart succeeded or failed, reading $SERVICE_RESULT to record
# the outcome as a node_exporter metric.
ExecStopPost={{ gitops_reconcile_metric_script }}
