# roles/autoupdate/templates/autoupdate.service.j2
[Unit]
Description=Unattended zypper {{ autoupdate_zypper_command }} and reboot
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/zypper --non-interactive {{ autoupdate_zypper_command }} --auto-agree-with-licenses
# 102/103 are zypper's "reboot/restart recommended" infos on a successful run,
# not failures: accept them so the reboot still fires. A real failure (a lower
# exit code) skips ExecStartPost, leaving the system un-rebooted for inspection.
SuccessExitStatus=102 103
ExecStartPost=/usr/bin/systemctl reboot
# Runs whether the update succeeded (then the reboot above) or failed, reading
# $SERVICE_RESULT to record the outcome as a node_exporter metric — a silently
# failing update otherwise leaves the host unpatched with no signal.
ExecStopPost=/usr/local/sbin/autoupdate-metric.sh
