# roles/autoupdate/templates/autoupdate.timer.j2
[Unit]
Description=Weekly unattended zypper dup

[Timer]
OnCalendar={{ autoupdate_oncalendar }}
RandomizedDelaySec=2h
# No Persistent=: if the box is off at the scheduled time, wait for next week
# rather than fire an unattended dup+reboot at an unpredictable post-boot moment.

[Install]
WantedBy=timers.target
