# roles/wordpress/files/wordpress-cron.service
[Unit]
Description=WordPress scheduled tasks
# wp-cron runs through the cli wrapper, which needs the stack up.
Requires=wordpress.service
After=wordpress.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wp cron event run --due-now
# Runs whether the run succeeded or failed, reading $SERVICE_RESULT to record the
# outcome as a node_exporter metric — a wp-cron timer that silently stops firing
# otherwise stalls every scheduled task unseen.
ExecStopPost=/usr/local/sbin/wp-cron-metric.sh
