# roles/arr/templates/app.container.j2
{% set arr_lsio = item.value.image.startswith('lscr.io/') %}
{% set arr_app_gid = arr_gid if item.key in arr_media_access else item.value.uid | default('') %}
{% set arr_data_suffix = ':ro' if arr_media_access[item.key] | default('rw') == 'ro' else '' %}
[Unit]
Description={{ item.key | capitalize }}
{% if item.value.data is defined %}
RequiresMountsFor={{ arr_data_root }}
{% endif %}
{% if item.value.netns is defined %}
{# Sharing a netns ties this app to its owner: restarts propagate via PartOf. #}
Requires={{ item.value.netns }}.service
After={{ item.value.netns }}.service
PartOf={{ item.value.netns }}.service
{% endif %}

[Container]
ContainerName={{ item.key }}
Image={{ item.value.image }}
{% if item.value.network | default('') == 'host' %}
Network=host
{% elif item.value.netns is defined %}
Network=container:{{ item.value.netns }}
{% else %}
Network=caddy.network
{% endif %}
{% if item.value.network_alias is defined %}
NetworkAlias={{ item.value.network_alias }}
{% endif %}
{% for cap in item.value.caps | default([]) %}
AddCapability={{ cap }}
{% endfor %}
{% for sysctl in item.value.sysctls | default([]) %}
Sysctl={{ sysctl }}
{% endfor %}
Environment=TZ={{ arr_timezone }}
{% if arr_lsio %}
Environment=PUID={{ item.value.uid }}
Environment=PGID={{ arr_app_gid }}
Environment=UMASK=002
{% elif item.value.uid is defined %}
User={{ item.value.uid }}:{{ arr_app_gid }}
{% endif %}
{% if item.value.device is defined %}
{# '-' marks the device optional, so GPU-less molecule containers still start. #}
AddDevice=-{{ item.value.device }}
{% endif %}
{# Host networking already exposes every port natively, so a publish is both #}
{# redundant and rejected by podman. #}
{% if item.value.publish | default(false) and item.value.network | default('') != 'host' %}
PublishPort={{ item.value.port }}:{{ item.value.port }}
{% endif %}
{# lscr.io images chown /config themselves; :U hands it to User= for the rest. #}
Volume={{ item.key }}-config:/config{{ '' if arr_lsio else ':U' }}
{% for bind in item.value.binds | default([]) %}
Volume={{ bind }}
{% endfor %}
{% for path in item.value.tmpfs | default([]) %}
Tmpfs={{ path }}
{% endfor %}
{% if item.value.data | default('') == 'root' %}
Volume={{ arr_data_root }}:/data{{ arr_data_suffix }}
{% elif item.value.data | default('') is mapping %}
{# A mapping pins each NAS subdir to an explicit container path (Plex's library #}
{# predates the /data convention), rather than the implicit /data/<name>. #}
{% for src, dest in item.value.data.items() %}
Volume={{ arr_data_root }}/{{ src }}:{{ dest }}{{ arr_data_suffix }}
{% endfor %}
{% else %}
{% for data in item.value.data | default([]) %}
Volume={{ arr_data_root }}/{{ data }}:/data/{{ data }}{{ arr_data_suffix }}
{% endfor %}
{% endif %}
{% if item.value.healthcheck is defined %}
HealthCmd={{ item.value.healthcheck.cmd }}
HealthInterval={{ item.value.healthcheck.interval | default('1m') }}
HealthTimeout={{ item.value.healthcheck.timeout | default('10s') }}
HealthStartPeriod={{ item.value.healthcheck.start_period | default('30s') }}
HealthRetries={{ item.value.healthcheck.retries | default(3) }}
{# A blackhole tunnel (no real wg config) is meant to stay down — never #}
{# restart-loop it; only act on failure once a real config is in place. #}
{% if item.key == 'wireguard' and arr_wireguard_conf | default('') | length == 0 %}
HealthOnFailure=none
{% else %}
HealthOnFailure={{ item.value.healthcheck.on_failure | default('none') }}
{% endif %}
{% endif %}

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
