# roles/arr/templates/app.container.j2
{% set arr_lsio = item.value.image.startswith('lscr.io/') %}
{% set arr_app_gid = arr_gid if item.key in arr_media_access else item.value.uid | default('') %}
{% set arr_data_suffix = ':ro' if arr_media_access[item.key] | default('rw') == 'ro' else '' %}
{% set arr_app_api_key = arr_api_keys.get(item.key, '') %}
[Unit]
Description={{ item.key | capitalize }}
{% if item.value.data is defined %}
RequiresMountsFor={{ arr_data_root }}
{% endif %}
{% if item.value.netns is defined %}
{# Sharing a netns ties this app to its owner: restarts propagate via PartOf. #}
Requires={{ item.value.netns }}.service
After={{ item.value.netns }}.service
PartOf={{ item.value.netns }}.service
{% endif %}

[Container]
ContainerName={{ item.key }}
Image={{ item.value.image }}
{% if item.value.network | default('') == 'host' %}
Network=host
{% elif item.value.netns is defined %}
Network=container:{{ item.value.netns }}
{% else %}
Network=caddy.network
{% endif %}
{# A netns owner advertises each enabled joiner that has a proxied port as a #}
{# network alias, so caddy resolves <app> to this container and reaches its #}
{# port inside the shared netns. A portless joiner (flaresolverr) needs none — #}
{# its netns mates reach it on localhost. #}
{% for peer in arr_apps | dict2items | selectattr('key', 'in', arr_enabled) if peer.value.netns | default('') == item.key and peer.value.port is defined %}
NetworkAlias={{ peer.key }}
{% endfor %}
{% for cap in item.value.caps | default([]) %}
AddCapability={{ cap }}
{% endfor %}
{% for sysctl in item.value.sysctls | default([]) %}
Sysctl={{ sysctl }}
{% endfor %}
Environment=TZ={{ arr_timezone }}
{% if arr_lsio %}
Environment=PUID={{ item.value.uid }}
Environment=PGID={{ arr_app_gid }}
Environment=UMASK=002
{% elif item.value.uid is defined %}
User={{ item.value.uid }}:{{ arr_app_gid }}
{% endif %}
{% if arr_app_api_key %}
{# API key from the repo, read from a 0600 EnvironmentFile so it stays out of #}
{# the world-readable unit; <APP>__AUTH__APIKEY overrides the config.xml key. #}
EnvironmentFile=/etc/arr/{{ item.key }}.env
{% endif %}
{% if item.key == 'transmission' and arr_transmission_username and arr_transmission_password %}
{# RPC auth creds from a 0600 EnvironmentFile (USER/PASS turn auth on in the LSIO #}
{# image and set the rpc user/password), kept out of the world-readable unit. #}
EnvironmentFile=/etc/arr/transmission.env
{% endif %}
{% if item.value.device is defined %}
{# '-' marks the device optional, so GPU-less molecule containers still start. #}
AddDevice=-{{ item.value.device }}
{% endif %}
{# lscr.io images chown /config themselves; :U hands it to User= for the rest. #}
Volume={{ item.key }}-config:/config{{ '' if arr_lsio else ':U' }}
{% for bind in item.value.binds | default([]) %}
Volume={{ bind }}
{% endfor %}
{% for path in item.value.tmpfs | default([]) %}
Tmpfs={{ path }}
{% endfor %}
{% if item.value.data | default('') == 'root' %}
Volume={{ arr_data_root }}:/data{{ arr_data_suffix }}
{% elif item.value.data | default('') is mapping %}
{# A mapping pins each NAS subdir to an explicit container path (Plex's library #}
{# predates the /data convention), rather than the implicit /data/<name>. #}
{% for src, dest in item.value.data.items() %}
Volume={{ arr_data_root }}/{{ src }}:{{ dest }}{{ arr_data_suffix }}
{% endfor %}
{% else %}
{% for data in item.value.data | default([]) %}
Volume={{ arr_data_root }}/{{ data }}:/data/{{ data }}{{ arr_data_suffix }}
{% endfor %}
{% endif %}
{% if item.value.healthcheck is defined %}
HealthCmd={{ item.value.healthcheck.cmd }}
HealthInterval={{ item.value.healthcheck.interval | default('1m') }}
HealthTimeout={{ item.value.healthcheck.timeout | default('10s') }}
HealthStartPeriod={{ item.value.healthcheck.start_period | default('30s') }}
HealthRetries={{ item.value.healthcheck.retries | default(3) }}
{# A blackhole tunnel (no real wg config) is meant to stay down — never #}
{# restart-loop it; only act on failure once a real config is in place. #}
{% if item.key == 'wireguard' and arr_wireguard_conf | default('') | length == 0 %}
HealthOnFailure=none
{% else %}
HealthOnFailure={{ item.value.healthcheck.on_failure | default('none') }}
{% endif %}
{% endif %}

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
